# playbooks/roles/jenkins_common/templates/etc/systemd/system/jenkins.service.j2
[Unit]
Description=Jenkins

[Service]
User=root
ExecStart=/usr/bin/java -jar /usr/share/jenkins/jenkins.war --httpPort={{ jenkins_port }}

[Install]
WantedBy=multi-user.target
